# etc/cron.d/zabbix_partitioning
# do partitioning 3 times a day 00:09 and 9:09, 18:09
9 */9 * * * root /usr/local/bin/zabbix_partitioning.sh 2>&1 >> /var/log/zabbix/zabbix_partitioning.log
